# etc/systemd/system/isucondition.go.service.d/override.conf
[Service]
LimitNOFILE=524288
StartLimitBurst=999
